# supervisor.d/phorge-sshd.conf
[program:phorge-sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config.phorge