# supervisor.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm8.4 -F